# ntpd.sv.conf
[program:ntpd]
command=/usr/sbin/ntpd -n
autostart=true
autorestart=true

